# libra-setup/DEBIAN/postrm
#!/bin/bash
set -e

SERVICE_USER="libra"

case "$1" in
    remove)
        # On remove, disable the service and remove files created by postinst
        if [ -d /run/systemd/system ]; then
            echo "Disabling libra service..."
            systemctl disable libra.service >/dev/null 2>&1 || true
        fi

        echo "Removing symbolic link for libphidget22..."
        rm -f /usr/lib/libphidget22.so

        echo "Reloading udev rules..."
        udevadm control --reload-rules
        udevadm trigger
    ;;

    purge)
        # On purge, remove the system user created by postinst.
        # This is called after 'remove' actions have already been performed.
        if getent passwd "$SERVICE_USER" >/dev/null; then
            echo "Removing system user '$SERVICE_USER'..."
            deluser --system "$SERVICE_USER" >/dev/null 2>&1 || true
        fi
    ;;

    upgrade|deconfigure|failed-upgrade|abort-install|abort-upgrade)
        # No action needed for these cases.
    ;;

    *)
        echo "postrm called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

exit 0